# backend/ocpp-server.service
# Systemd сервис для OCPP сервера
# Сохранить как /etc/systemd/system/evpower-ocpp.service

[Unit]
Description=EvPower OCPP Server - Pure FastAPI
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=exec
User=evpower
Group=evpower
WorkingDirectory=/opt/evpower-ocpp/backend
Environment=PATH=/opt/evpower-ocpp/backend/venv/bin
EnvironmentFile=/opt/evpower-ocpp/backend/.env
ExecStart=/opt/evpower-ocpp/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8180 --workers 1
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

# Безопасность
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/evpower-ocpp/backend
ReadWritePaths=/var/log/evpower-ocpp
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=evpower-ocpp

[Install]
WantedBy=multi-user.target 